# res/gitlab-workhorse/systemd/gitlab-workhorse.socket
[Socket]
ListenStream=/var/run/gitlab-workhorse.socket
